# deploy/fix-life-celery-beat.service
[Unit]
Description=Fix Life Celery Beat
After=network-online.target fix-life-celery-worker.service

[Service]
Type=simple
User=josie
Group=josie
WorkingDirectory=/opt/fix-life/backend
Environment="PATH=/opt/fix-life/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/fix-life/backend/.venv/bin/celery -A app.core.celery beat --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
